# deploy/watchpi-deploy.service
[Unit]
Description=WatchPi auto-deploy (pull latest from GitHub)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/srv/apps/watchpi/deploy/autodeploy.sh
Nice=10
